# build/core/containers/provisioner/cookbooks/default/templates/etc/supervisor/conf.d/gush.conf.erb
[program:gush]
#command=bash -c 'source /usr/share/rvm/scripts/rvm && /projects/ansible/provisioner/run_gush.sh <%= node['attributes'].fetch('env') %>'
command=bash -c 'source /usr/share/rvm/scripts/rvm && /projects/ansible/gush/run.sh <%= node['attributes'].fetch('env') %>'
autostart=true
autorestart=true
stderr_logfile=/var/log/gush/out.log
stdout_logfile=/var/log/gush/out.log
