# build/core/cookbooks/common_cookbooks/ex2/files/supervisor/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld --defaults-file=/etc/mysql-default/my.cnf --basedir=/usr --pid-file=/var/run/mysql-default/mysqld.pid
autostart=true
autorestart=true
user=mysql
